# roles/core/backup/templates/monitoring-health.service.j2
[Unit]
Description=Monitoring Stack Health Check Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/health-check.sh
User=root
Group=root

[Install]
WantedBy=multi-user.target
